# arrow.service
[Unit]
Description=Arrow AI Assistant
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspaces/Arrow
ExecStart=/bin/bash /workspaces/Arrow/run_arrow.sh
Restart=always
RestartSec=5
StandardOutput=append:/workspaces/Arrow/arrow_data/logs/arrow.log
StandardError=append:/workspaces/Arrow/arrow_data/logs/arrow.err
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
